# scripts/sample-report-deploy.path
[Unit]
Description=Watch for deploy trigger file

[Path]
PathExists=/home/anurag/projects/sample_report_handler/.deploy-trigger
Unit=sample-report-deploy.service

[Install]
WantedBy=multi-user.target
